# scripts/orderbook_collector.service
[Unit]
Description=Trade Platform — Orderbook DOM Collector
Documentation=https://github.com/yourrepo/trade_platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/vscode_projects/trade_platform

# Adjust the python path if your venv is named differently
ExecStart=/root/vscode_projects/trade_platform/venv/bin/python \
    /root/vscode_projects/trade_platform/scripts/collect_orderbook.py

# Restart on any failure (network drops, API errors, etc.)
Restart=always
RestartSec=10

# Optional: override defaults via environment variables
# Environment=DOM_INTERVAL=60
# Environment=DOM_DEPTH_LIMIT=50
# Environment=DB_PATH=/root/vscode_projects/trade_platform/trade_platform.db

# Capture output to systemd journal (view with: journalctl -u orderbook-collector -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orderbook-collector

[Install]
WantedBy=multi-user.target
